# services/plex/plex.service
[Unit]
Description=A service that runs Plex
After=network.target mnt-files.mount docker.service
Requires=mnt-files.mount docker.service

[Service]
ExecStart=/bin/bash -c 'exec $DROID_HOME/services/plex/start.sh'
ExecStop=/bin/bash -c 'docker-compose --file $DROID_HOME/services/plex/docker-compose.yml \
                            --project-name plex stop'
Restart=on-failure

[Install]
WantedBy=default.target
